# packaging/systemd/privatetv-schedule.timer
[Unit]
Description=Hourly PrivateTV schedule maintenance

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
